# ansible/roles/common/files/69-hdparm.rules.tmpl
ACTION=="add|change", KERNEL=="sd[a-z]", ATTRS{queue/rotational}=="1", RUN+="{{ hdparm_path }} -B 20 -S 5/dev/%k"